# docker/web/build.ini
[program:weekly-web]
directory = /data/www/web
command = sh build.sh
